# nix/nix-tarball/store/2qzfvsqb9afhb73cc3yfg8hk2xpxcy47-nix-2.16.0pre20230512_dirty-riscv64-unknown-linux-gnu/etc/init/nix-daemon.conf
description "Nix Daemon"
start on filesystem
stop on shutdown
respawn
exec /nix/store/2qzfvsqb9afhb73cc3yfg8hk2xpxcy47-nix-2.16.0pre20230512_dirty-riscv64-unknown-linux-gnu/bin/nix-daemon --daemon
